# roles/nextcloud/templates/nextcloud-cron.service
[Unit]
Description=Nextcloud cron.php job

[Service]
User=apache
ExecStart=/usr/bin/php -f /var/www/html/cron.php

[Install]
WantedBy=basic.target
